# scripts/2021/donaraltna@.service
[Unit]
Description=donaraltna

[Service]
ExecStart=/usr/local/bin/donaraltna %i
WorkingDirectory=/root/donaraltna/

