# scripts/2021/donardupopt@.service
[Unit]
Description=donardupopt

[Service]
ExecStart=/usr/local/bin/donardupopt %i
WorkingDirectory=/root/donardupopt/

