# scripts/2021/donardupna@.service
[Unit]
Description=donardupna

[Service]
ExecStart=/usr/local/bin/donardupna %i
WorkingDirectory=/root/donardupna/

