# scripts/2021/donaraltopt@.service
[Unit]
Description=donaraltopt

[Service]
ExecStart=/usr/local/bin/donaraltopt %i
WorkingDirectory=/root/donaraltopt/

