# scripts/2021/donaralt@.service
[Unit]
Description=donaralt

[Service]
ExecStart=/usr/local/bin/donaralt %i
WorkingDirectory=/root/donaralt/

