# scripts/2021/donardup@.service
[Unit]
Description=donardup

[Service]
ExecStart=/usr/local/bin/donardup %i
WorkingDirectory=/root/donardup/

