# scripts/2021/torexit@.service
[Unit]
Description=torexit

[Service]
ExecStart=/usr/local/bin/torexit %i
WorkingDirectory=/root/torexit/
