# scripts/2021/torexithelper.service
[Unit]
Description=torexit

[Service]
ExecStart=/usr/local/bin/torecho -s 172 -p 443

