# scripts/2021/torexitna@.service
[Unit]
Description=torexitna

[Service]
ExecStart=/usr/local/bin/torexitna %i
WorkingDirectory=/root/torexitna/
