# scripts/2021/torfone@.service
[Unit]
Description=torfone

[Service]
ExecStart=/usr/local/bin/torfone %i
WorkingDirectory=/root/torfone/

