# scripts/2021/torfoneopt@.service
[Unit]
Description=torfoneopt

[Service]
ExecStart=/usr/local/bin/torfoneopt %i
WorkingDirectory=/root/torfoneopt/

