# scripts/2021/torfonena@.service
[Unit]
Description=torfonena

[Service]
ExecStart=/usr/local/bin/torfonena %i
WorkingDirectory=/root/torfonena/

