# scripts/2021/torhsna@.service
[Unit]
Description=torhsna

[Service]
ExecStart=/usr/local/bin/torhsna %i
WorkingDirectory=/root/torhsna/

