# scripts/2021/torhsopt@.service
[Unit]
Description=torhsopt

[Service]
ExecStart=/usr/local/bin/torhsopt %i
WorkingDirectory=/root/torhsopt/

