# scripts/2021/torhs@.service
[Unit]
Description=torhs

[Service]
ExecStart=/usr/local/bin/torhs %i
WorkingDirectory=/root/torhs/

